# ping-monitor.service
[Unit]
Description=Pi Ping Monitor — automatic network failover service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=@@REAL_USER@@
Environment=CONFIG_FILE=/etc/ping-monitor/config.env
LogsDirectory=ping-monitor
StateDirectory=ping-monitor
ExecStart=/usr/local/bin/ping-monitor.sh
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
